# playmock/roles/container_gitea/files/backup/cron.d/backup_container_gitea
# container_gitea backup schedule
05 03 * * * root /opt/containers/backup/gitea_backup.sh >/dev/null 2>&1
