# playmock/roles/container/files/utils/cron.d/container_utils
# container utilities schedule
55 03 * * * root /opt/containers/utils/container_cleanup.sh >/dev/null 2>&1
