# playmock/roles/container_redmine/files/backup/cron.d/backup_container_redmine
# container_redmine backup schedule
05 02 * * * root /opt/containers/backup/redmine_backup.sh >/dev/null 2>&1
