# playmock/roles/idhub_oidc/templates/idhub-oidc-green.service.j2
[Unit]
Description=idhub-oidc-green
After=syslog.target

[Service]
EnvironmentFile=/usr/local/app/idhub-oidc/oidc-green.conf
WorkingDirectory = /usr/local/app/idhub-oidc
ExecStart=/usr/local/app/idhub-oidc/start.sh
SuccessExitStatus=143
User=gooscp
Group=gooscp

[Install]
WantedBy=multi-user.target
