# playmock/roles/container_jenkins/files/backup/cron.d/backup_container_jenkins
# container_jenkins backup schedule
05 01 * * * root /opt/containers/backup/jenkins_backup.sh >/dev/null 2>&1
