# playmock/roles/common/templates/usr/lib/systemd/system/dnf-automatic-rpm-collection.service.j2
[Unit]
Description=rpm collection service

[Service]
#User=root (デフォルトのままなので省略)
Type=oneshot
ExecStart=/usr/local/bin/dnf-auto-rpm-collection.sh

[Install]
WantedBy=multi-user.target
